# lager-backend.service
[Unit]
Description=Lagerverwaltung Backend (Express + SQLite)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/cartoon-lager-blick/backend
ExecStart=/usr/bin/npm run dev
Restart=always
User=www-data
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
